# docs/example/separator.socket
[Unit]
Description=separator socket

[Socket]
ListenStream=/run/separator.sock

[Install]
WantedBy=sockets.target